# packaging/cbhm.service
[Unit]
Description=Clipboard helper module
After=xorg.target

[Service]
Type=simple
ExecStart=/usr/bin/cbhm

[Install]
WantedBy=core-efl.target
